# meta-pike/recipes-core/underclock/underclock/underclock.service
[Unit]
Description=Underclock CPU to save battery life

[Service]
Type=oneshot
ExecStart=/usr/bin/underclock

[Install]
WantedBy=basic.target
